# cron.daily/aggdb
#!/bin/bash

/usr/bin/python /opt/IPhistdb/aggregate.py

