# cron.daily/iphistdb-maintenance
#!/bin/bash

/usr/bin/python /opt/IPhistdb/maintenance.py

