# cron.d/iphistdb-aggregate
0 0 * * *     root   /usr/bin/python /opt/IPhistdb/aggregate.py

